# cascadya-edge-os-images/mender/rootfs_overlay/etc/systemd/system/mender-updated.service.d/10-requires-data.conf
[Unit]
RequiresMountsFor=/data
After=data.mount mender-data-dir.service
